# templates/default/upstart/adam.conf.erb
description "Adam meta-service"

start on runlevel [2345]
stop on runlevel [016]
