# templates/default/sysvinit/component.erb
#!/bin/sh
#
### BEGIN INIT INFO
# Provides:          <%= @component_name %>
# Required-Start:    $remote_fs $syslog $network
# Required-Stop:     $remote_fs $syslog $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start this Adam component
# Description:       Component of Adam
#### END INIT INFO

. /lib/lsb/init-functions

PATH=/opt/rbenv/shims:/opt/rbenv/bin:/usr/local/sbin:/usr/local/bin:/usr/bin:/usr/sbin:/sbin:/bin
export PATH

DESC="Adam <%= @component_name %> component"
NAME=adam-<%= @component_name %>
ADAM_BASEDIR="<%= @base_directory %>"
DAEMON=${ADAM_BASEDIR}/bin/<%= @component_name %>
DAEMON_ARGS=""
LOG_FILE=/var/log/adam/<%= @component_name %>.log
ADAM_ENV="/etc/adam/environment"
PIDFILE="/var/run/adam/<%= @component_name %>.pid"
USER=adam

start() {
	source $ADAM_ENV
	start-stop-daemon --start --quiet --background --make-pidfile \
		--pidfile $PIDFILE --exec /bin/bash --name $NAME --user $USER \
		--chdir $ADAM_BASEDIR -- -c "exec $DAEMON $DAEMON_ARGS >> $LOG_FILE 2>&1" \
		|| return 2
}

stop() {
	start-stop-daemon --stop --quiet \
		--pidfile $PIDFILE \
		--retry=TERM/30/KILL/5
}

case "$1" in
start)
	log_daemon_msg "Starting Adam Component <%= @component_name %>" || true
	if start; then
		log_end_msg 0 || true
	else
		log_end_msg 1 || true
	fi
	;;
stop)
	log_daemon_msg "Stopping Adam Component <%= @component_name %>" || true
	if stop; then
		log_end_msg 0 || true
	else
		log_end_msg 1 || true
	fi
;;
restart)
	log_daemon_msg "Restarting Adam Component <%= @component_name %>" || true
	if stop && start; then
		log_end_msg 0 || true
	else
		log_end_msg 1 || true
	fi
	;;
force-reload)
	log_action_msg "force-reload not supported"
	exit 1
	;;
status)
	status_of_proc -p $PIDFILE "$DAEMON" "$NAME" && exit 0 || exit $?
	;;
esac

exit 0
